# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable energieip-swh200-core.service
systemctl restart energieip-swh200-core.service
