# scripts/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable core.service
systemctl start core.service
